# src/dockers/php/supervisord/conf/crontab.conf
### to schedule run, uncomment all
[program:crontab]
command=/usr/bin/crontab -f
autostart=true
autorestart=false
startsecs = 0
startretries = 1
user=root



